# rootfs/etc/services.d/ssh/run
#!/bin/zsh
su git -c '/usr/sbin/sshd -D -E /var/log/gitd/sshd.log &'  
tail -f /var/log/gitd/sshd.log 