# rootfs/etc/services.d/gitd-log/run
#!/bin/zsh
tail -f /var/log/gitd/gitd.log 
